# src/pi-powerbutton.service
[Unit]
Description=Power Button

[Service]
ExecStart=/usr/bin/pi-powerbutton
Restart=always

[Install]
WantedBy=default.target

